# chronos-backup.timer
[Unit]
Description=Chronos Lab Daily Backup Timer
After=network.target

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=chronos-backup.service

[Install]
WantedBy=timers.target
