# chronos-backup.service
[Unit]
Description=Chronos Lab Daily Backup
After=network.target

[Service]
Type=oneshot
ExecStart=/home/claworc/.openclaw/workspace/scripts/backup-to-github.sh
User=claworc
WorkingDirectory=/home/claworc/.openclaw/workspace
StandardOutput=append:/var/log/chronos-backup.log
StandardError=append:/var/log/chronos-backup.log

[Install]
WantedBy=multi-user.target
